# TestPythonService.service
[Unit]
Description=A test python service
After=multi-user.target

[Service]
Type=idle
Environment="DISPLAY=:0"
ExecStart=/usr/bin/python3 /usr/local/bin/start.py

[Install]
WantedBy=multi-user.target
